# crontab
0 7 * * * /app/entrypoint.sh >> /var/log/menucron.log 2>&1
