# crontab
1 0 * * * /app/cron.sh >> /app/log/menucron.log 2>&1
1 12 * * * /app/cron.sh >> /app/log/menucron.log 2>&1
